# deployment_files/todo_app.service
[Unit]
description=Todo app service

[Service]
ExecStart=/opt/todoapp/run_server.sh
WorkingDirectory=/opt/todoapp
User=ec2-user
